# vm/polymarket.service
[Unit]
Description=Polymarket WebSocket Ingest
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/var/data/polymarket
[Service]
Type=simple
EnvironmentFile=/etc/polymarket/polymarket.env
User=root
ExecStartPre=/usr/bin/env bash -lc 'test -x "$APP_DIR/.venv/bin/python"'
ExecStartPre=/usr/bin/env bash -lc 'test -f "$APP_DIR/$ENTRYPOINT"'
ExecStart=/usr/bin/env bash -lc 'cd "$APP_DIR" && exec "$APP_DIR/.venv/bin/python" "$APP_DIR/$ENTRYPOINT"'
Restart=always
RestartSec=3
LimitNOFILE=262144
[Install]
WantedBy=multi-user.target
