# vm/clean/polymarket-clean.service
[Unit]
Description=Polymarket clean local parquet cache
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/polymarket/polymarket.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/local/bin/clean_old_local.sh
